# config/supervisor-salt.conf
[program:salt-master]
command=/usr/bin/salt-master --log-level=%(ENV_LOG_LEVEL)s
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
